# deploy-scripts/webapi.service
[Unit]
Description=.NET Dice Roll Web API

[Service]
WorkingDirectory=/var/www/
ExecStart=/usr/bin/dotnet /var/www/WebAdvert.Web.dll